# centos/init_supervisord/etc/supervisord.d/syslog-ng.conf
[program:syslog-ng]
command=/usr/sbin/syslog-ng --foreground  
redirect_stderr=true
autostart=true
stdout_syslog=true

